# varnish/etc/systemd/system/varnish@mbharata.service
[Unit]
Description=Varnish Cache, a high-performance HTTP accelerator
Documentation=https://www.varnish-cache.org/docs/ man:varnishd

[Service]
Type=simple

# Maximum number of open files (for ulimit -n)
LimitNOFILE=131072

# Locked shared memory - should suffice to lock the shared memory log
# (varnishd -l argument)
# Default log size is 80MB vsl + 1M vsm + header -> 82MB
# unit is bytes
LimitMEMLOCK=85983232

#ExecStartPre=sleep 3 ; sudo rm /var/lib/varnish/sin1//_.pid ; true
ExecStart=/usr/sbin/varnishd \
          -P /run/varnish@mbharata.pid \
          -j unix,user=vcache \
          -i 10812 \
          -F \
          -a :10812 \
          -f /etc/varnish/mbharata.vcl \
          -S /etc/varnish/secret \
          -s default=file,/var/varnish/mbharata,16G \
          -s mem=malloc,16m
#          -T localhost:6082 \

#ExecReload=/usr/share/varnish/varnishreload
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true

[Install]
WantedBy=multi-user.target
